# startup.cron
55,56,57,58,59,0 18,19,20 * * * echo "Hello, there"
