# startup.cron
00 * * * * /home/ec2-user/cw2t/cw2tstartup.sh
