# startup.cron
MAILTO="info@yearstretch.com"
45 18 * * * cw2tstartup.sh
