# startup.cron
SHELL=/bin/sh
15 * * * * ../node_modules/forever/bin/forever list
