# startup.cron
55,56,57,58,59 * * * * echo "Hello, there"
