# startup.cron
SHELL=/bin/sh
05 * * * * ../node_modules/forever/bin/forever list
