# startup.cron
10 10 * * * ../node_modules/forever/bin/forever list
